# dotfiles/.config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Unit]
After=dbus.service
Requires=dbus.service
